# rio_api/aventura/etc/systemd/system/rioos_api.service.erb
[Unit]
Description=Rio/OS <%= @package[:package] %>
After=network.target

[Service]
ExecStart=<%= @basic[:home]+"/"+@package[:package]+"/bin" %>/rioos-apiserver start
KillMode=process

[Install]
WantedBy=multi-user.target
