# rio_prometheus/aventura/etc/systemd/system/rioos_prometheus.service.erb
[Unit]
Description=<%= @package[:package] %>
After=network.target

[Service]
ExecStart=<%= @basic[:home] %>/<%= @package[:package] %>/bin/prometheus --config.file=<%= @basic[:libhome] %>/<%= @package[:package] %>/prometheus_rioos.yml
KillMode=process

[Install]
WantedBy=multi-user.target
