# rio_nilavu/xenial/etc/systemd/system/rioos-ui.service.erb
[Unit]
Description=Rio/OS Console
After=network.target

[Service]
WorkingDirectory=/var/www/<%= @package[:package] %>/apps/nilavu
ExecStart=/usr/local/bin/ember serve
KillMode=process

[Install]
WantedBy=multi-user.target
