# rio_storlet/centos7/etc/systemd/system/rioos.service.erb
[Unit]
Description=<%= @package[:package] %> core engine - Omni scheduler.
After=network.target
After=runlevel2.target
After=runlevel3.target
After=runlevel4.target
After=runlevel5.target
[Service]
EnvironmentFile=<%= @basic[:libhome] %>/env.sh
ExecStart=<%= @basic[:home] %>/<%= @package[:package] %>/bin/<%= @package[:package] %> -v start
KillMode=process
