# rio_vnc/centos7/etc/systemd/system/rioosvnc.service.erb
[Unit]
Description=Verticevnc Server
After=network.target
After=runlevel2.target
After=runlevel3.target
After=runlevel4.target
After=runlevel5.target

[Service]
ExecStart=/usr/bin/node <%= @basic[:home]+"/"+@package[:package] %>/vnc_monitor.js
KillMode=process

[Install]
WantedBy=multi-user.target
