# rio_gulpd/centos7/etc/systemd/system/verticegulpd.service.erb
[Unit]
Description=<%= @package[:package] %> Agent
After=local-fs.target network.target
Requires=network.target
Wants=local-fs.target

[Service]
User=root
Group=root
EnvironmentFile=<%= @basic[:libhome] %>/env.sh
ExecStart=<%= @basic[:home] %>/<%= @package[:package] %>/bin/gulpd  start --config <%= @basic[:libhome] %>/gulp/gulpd.conf
KillMode=process
Restart=on-failure
RestartSec=42s

[Install]
WantedBy=multi-user.target
