# rio_scheduler/aventura/etc/systemd/system/rioos_scheduler.service.erb
[Unit]
Description=<%= @package[:package] %>
After=network.target

[Service]
ExecStart=<%= @basic[:home]+"/"+@package[:package]+"/bin" %>/rioos-scheduler --v=4  --api-server=http://localhost:9636
KillMode=process

[Install]
WantedBy=multi-user.target
