# rio_scheduler/aventura/etc/systemd/system/rioos_scheduler.service.erb
[Unit]
Description=<%= @package[:package] %>
After=network.target

[Service]
ExecStart=<%= @basic[:home]+"/"+@package[:package]+"/bin" %>/rioos-scheduler --v=4 --leader-elect=false --api-server=https://localhost:7443  --watch-server=https://localhost:8443 --service-account-private-key-file=<%= @basic[:libhome] %>/config/service-account.key --use-service-account-credentials --rioconfig=<%= @basic[:libhome] %>/config/scheduler.rioconfig
KillMode=process

[Install]
WantedBy=multi-user.target
