# rio_fluentbit/aventura/etc/systemd/system/rioos-fluentbit.service.erb
[Unit]
Description=Fluent-bit Log server
After=network.target

[Service]
ExecStart=/usr/local/bin/fluent-bit -c <%= @basic[:home]+"/"+@package[:package]+"/conf/fluent-bit.conf" %>
KillMode=process

[Install]
WantedBy=multi-user.target
